# contrib/systemd/batt_log.service
[Unit]
Description=batt_log service

[Service]
ExecStart=/usr/local/bin/batt_log-daemon
Restart=always

[Install]
WantedBy=multi-user.target

