# service/services/openrc/batt_log
#!/sbin/openrc-run

description="Script to execute batt_log"
command="/usr/local/bin/batt_log"
command_args=""
